# src/etc/init.d/mudimodem-collectd
#!/bin/sh /etc/rc.common
# procd service for the MudiModem background history collector (spec §6).
# Polls ubus (not AT) and appends samples to /tmp/mudimodem. Respawned on crash.

START=95
STOP=10
USE_PROCD=1
PROG=/usr/sbin/mudimodem-collectd

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/python3 "$PROG"
	procd_set_param respawn
	procd_set_param stderr 1
	procd_close_instance
}
